# files/etc/hotplug.d/iface/30-usb
#!/bin/bash

[[ "$INTERFACE" == "usb" ]] || exit
[[ "$ACTION" == "ifup" ]] || exit

sleep 3
DEF_ROUTE=$(route | grep default | awk '{print $2}')
/usr/bin/pineapple/usb0_ping $DEF_ROUTE
