# files/etc/hotplug.d/usb/30-sd
#!/bin/bash


[[ $PRODUCT == "5e3/745/933" ]] && {

  [[ $ACTION == "add" ]] && {
    echo "add" > /tmp/sd_action
  } || {
    echo "remove" > /tmp/sd_action
  }

}

